# scripts/manjaro/etc/systemd/system/open-fprintd-resume.service.d/override.conf
[Service]
ExecStart=
ExecStart=systemctl restart python3-validity.service open-fprintd.service
